# deploy/systemd/sport-prediction-frontend.service
[Unit]
Description=Sport Prediction Frontend Static Server
After=network.target

[Service]
Type=simple
User=sportapp
Group=sportapp
WorkingDirectory=/opt/sport-prediction/current/frontend
ExecStart=/opt/sport-prediction/current/frontend/node_modules/.bin/serve -s dist -l tcp://0.0.0.0:8101
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
StandardOutput=append:/var/log/sport-prediction/frontend.log
StandardError=append:/var/log/sport-prediction/frontend-error.log

[Install]
WantedBy=multi-user.target
