# deploy/systemd/sport-prediction-discovery.timer
[Unit]
Description=Run Sport Prediction Stage 1 Discovery every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
